# tasmota-ccc/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Custom Add-on: Tasmota proxy and update
# Runs the Nginx daemon
# ==============================================================================

source /var/lib/tasmota-ccc/functions.sh

# Generate Ingress configuration
interface="$(bashio::addon.ip_address)"
port="$(bashio::addon.ingress_port)"
bashio::log.info "Running on ${interface}:${port}"
bashio::var.json \
    interface "${interface}" \
    port "^${port}" \
    | tempio \
        -template /etc/nginx/templates/nginx.gtpl \
        -out /etc/nginx/nginx.conf

# Generate Proxy configurations
echo "<body><table align="center">" > /var/www/tasmota-ccc/devices.html
get_tasmota_devices | jq -rcM "sort_by(.name) | .[]" | while read device; do
    # id, available, name, model, version, area, url, mac
    eval "declare -A device_attributes=($(bashio::jq "${device}" 'to_entries | .[] | @sh "[\(.key)]=\(.value)"'))"

    bashio::var.json \
        id "${device_attributes[id]}" \
        url "${device_attributes[url]}" \
        | tempio \
            -template /etc/nginx/templates/device_proxy.gtpl \
            -out /etc/nginx/locations/device_proxy_"${device_attributes[id]}".conf

    bashio::var.json \
        id "${device_attributes[id]}" \
        name "${device_attributes[name]}" \
        area "${device_attributes[area]}" \
        model "${device_attributes[model]}" \
        version "${device_attributes[version]}" \
        | tempio \
            -template /var/www/tasmota-ccc/templates/device.gtpl \
            -out /tmp/device.html
    cat /tmp/device.html >> /var/www/tasmota-ccc/devices.html
done
echo "</table></body>" >> /var/www/tasmota-ccc/devices.html

bashio::log.info "Starting NGinx..."
exec nginx -g "daemon off;"
